# deploy/xiaobo-agent.service
[Unit]
Description=小柏 Agent - 个人数字伙伴
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/xiaobo-agent
ExecStart=/usr/bin/python3 main.py --daemon
Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=60

# 环境变量
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xiaobo-agent

# 安全加固
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/root/.xiaobo-agent
ReadOnlyPaths=/root/xiaobo-agent

[Install]
WantedBy=multi-user.target
